# deploy/systemd/scrobblebox-oscilloscope.service
[Unit]
Description=ScrobbleBox Oscilloscope
After=network-online.target scrobblebox-core.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/scrobblebox
EnvironmentFile=/opt/scrobblebox/.env
ExecStart=/opt/scrobblebox/.venv/bin/python -m scrobblebox.oscilloscope.service
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
